# crontab
#!/bin/bash

*/2 * * * * /usr/local/bin/php /app/man invoice:cancel >/dev/null 2>&1